# backend/deploy/bince-bot.service
[Unit]
Description=Bince 24/7 Hybrid Trading Bot
After=network.target

[Service]
Type=simple
User=bince
WorkingDirectory=/opt/bince
ExecStart=/opt/bince/.venv/bin/python /opt/bince/hybrid_trader.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

# NOTE: This unit runs the bot directly as an always-on service — the most
# robust option for production (survives reboots, auto-restarts on crash).
#
# Alternatively, the API can start/stop the bot on demand via /bot/start and
# /bot/stop (see bot_service.py). Use ONE approach, not both, to avoid two
# instances trading the same pairs. If you enable this systemd unit, hide the
# Bot Control start/stop buttons or point them at `systemctl` instead.
