# backend/deploy/bince-worker.service
[Unit]
Description=Bince Optimization Worker (RQ)
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
User=bince
WorkingDirectory=/opt/bince/backend
Environment=REDIS_URL=redis://localhost:6379/0
Environment=PYTHONPATH=/opt/bince:/opt/bince/backend
ExecStart=/opt/bince/.venv/bin/rq worker optimizer --url redis://localhost:6379/0
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
